# rootfs/etc/s6-overlay/s6-rc.d/xorg-server/finish
#!/command/execlineb -P

define SOCKET_DIR /tmp/.X11-unix

if { test -d ${SOCKET_DIR} } rm -rf ${SOCKET_DIR}
