# rootfs/etc/s6-overlay/s6-rc.d/syslogd/run
#!/command/execlineb -P

exec syslogd -n -t -O -
