# rootfs/etc/s6-overlay/s6-rc.d/neko/run
#!/command/execlineb -P

with-contenv

define PULSE_SERVER unix:/tmp/pulseaudio.socket
define NEKO_SERVER_BIND :8080
define NEKO_PLUGINS_ENABLED true
define NEKO_PLUGINS_DIR /etc/neko/plugins/

importas -i -S X_NEKO_DISPLAY
importas -i -S X_NEKO_USER
importas -i -S X_NEKO_USER_HOME
importas -i -S X_NEKO_RUNTIME_DIR

foreground { mkdir -p -m 1777 ${SOCKET_DIR} }

export USER ${X_NEKO_USER}
export HOME ${X_NEKO_USER_HOME}
export DISPLAY ${X_NEKO_DISPLAY}
export PULSE_SERVER ${PULSE_SERVER}
export NEKO_SERVER_BIND ${NEKO_SERVER_BIND}
export NEKO_PLUGINS_ENABLED ${NEKO_PLUGINS_ENABLED}
export NEKO_PLUGINS_DIR ${NEKO_PLUGINS_DIR}
export XDG_RUNTIME_DIR ${X_NEKO_RUNTIME_DIR}

s6-setuidgid ${X_NEKO_USER}

exec neko serve --server.static "/var/www"
