# rootfs/etc/s6-overlay/s6-rc.d/debug-user-init/up
#!/command/execlineb -P

with-contenv

importas -S X_NEKO_USER
importas -S X_NEKO_GROUP
importas -S X_NEKO_UID
importas -S X_NEKO_GID
importas -S X_NEKO_USER_HOME
importas -S X_NEKO_RUNTIME_DIR
importas -S X_NEKO_DISPLAY

foreground { echo X_NEKO_USER = ${X_NEKO_USER} }
foreground { echo X_NEKO_GROUP = ${X_NEKO_GROUP} }
foreground { echo X_NEKO_UID = ${X_NEKO_UID} }
foreground { echo X_NEKO_GID = ${X_NEKO_GID} }
foreground { echo X_NEKO_USER_HOME = ${X_NEKO_USER_HOME} }
foreground { echo X_NEKO_RUNTIME_DIR = ${X_NEKO_RUNTIME_DIR} }
foreground { echo X_NEKO_DISPLAY = ${X_NEKO_DISPLAY} }
