# rootfs/etc/s6-overlay/s6-rc.d/sxhkd/run
#!/command/execlineb -P

with-contenv

importas -i -S X_NEKO_DISPLAY
importas -i -S X_NEKO_USER
importas -i -S X_NEKO_USER_HOME
importas -i -S X_NEKO_USER_SHELL

export DISPLAY ${X_NEKO_DISPLAY}
export HOME ${X_NEKO_USER_HOME}
export SHELL ${X_NEKO_USER_SHELL}
export USER ${X_NEKO_USER}
export XDG_DATA_DIRS /usr/local/share:/usr/share

s6-setuidgid ${X_NEKO_USER}
cd ${X_NEKO_USER_HOME}
exec /usr/bin/sxhkd
